# callme.ini
[program:callme]
command=php -f ~/Listener/callIn.php
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stderr_logfile=~/Listener/daemon.log
stdout_logfile=~/Listener/daemon.log

